# app/conf/supervisord.conf
[unix_http_server]
file=/mydata/run/supervisor.sock   ; (the path to the socket file)
chmod=0777
[supervisord]
logfile=/mydata/run/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
stderr_logfile=/mydata/run/error.log;
pidfile=/mydata/run/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false              ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)
childlogdir=/mydata/run/ ;
#user=root;

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///mydata/run/supervisor.sock ; use a unix:// URL  for a unix socket




[program:push_followee_service]
process_name = push_followee_service-%(process_num)s
command = /mydata/python2.7/bin/python /mydata/python/live_video/background/push_followee_service.py runfcgi daemonize=false method=prefork host=127.0.0.1 port=%(process_num)s
stdout_logfile = /mydata/logs/push_followee_service.log
autostart = true
autorestart = true
numprocs_start = 1
numprocs= 8

[program:push_friends_service]
process_name = push_friends_service-%(process_num)s
command = /mydata/python2.7/bin/python /mydata/python/live_video/background/push_friends_service.py runfcgi daemonize=false method=prefork host=127.0.0.1 port=%(process_num)s
stdout_logfile = /mydata/logs/push_friends_service.log
autostart = true
autorestart = true
numprocs_start = 1
numprocs= 8

[program:push_reportadmin_service]
process_name = push_reportadmin_service-%(process_num)s
command=/mydata/python2.7/bin/python /mydata/python/live_video/background/push_reportadmin_service.py runfcgi daemonize=false method=prefork host=127.0.0.1 port=%(process_num)s
process_name=%(program_name)s_%(process_num)02d
stdout_logfile = /mydata/logs/push_reportadmin_service.log
autostart = true
autorestart = true
numprocs_start = 1
numprocs= 8